# contrib/systemd/rbl-rbldns@spamip.service
.include /lib/systemd/system/rbl-rbldns@.service

[Unit]
Description = RBLDNS Exporter for "Spam IP" list

[Install]
WantedBy = multi-user.target
